# root/etc/services.d/webserver/run
#!/usr/bin/with-contenv bash

COMMAND="python3 -m http.server"

cd /map

exec s6-setuidgid abc $COMMAND
